# modules/rootfs/etc/profile.d/bash_aliases.sh
#!/bin/sh

# detailed file and directory listing
alias ll='ls -alih --color=always'
alias ls='ls --color=always'
alias rm='trash -v'
alias cp='cp -i'
alias mv='mv -i'
