# modules/rootfs/etc/systemd/system/update-mirrors.service
[Unit]
Description=Update pacman mirrorlist

[Service]
Type=simple
ExecStart=/usr/bin/rate-mirrors --allow-root --save /etc/pacman.d/mirrorlist arch
