# modules/rootfs/etc/systemd/system/update-mirrors.timer
[Unit]
Description=Update pacman mirrorlist

[Timer]
OnCalendar=daily
OnBootSec=1min

[Install]
WantedBy=timers.target

