# modules/rootfs/home/.config/systemd/user/rclone-gdrive.path
[Path]
Unit=rclone-gdrive.service
PathChanged=%h/GoogleDrive
