# modules/rootfs/home/.config/systemd/user/rclone-onedrive.path
[Path]
Unit=rclone-onedrive.service
PathChanged=%h/OneDrive
