# modules/rootfs/home/.config/systemd/user/rclone-gdrive-sync.timer
[Unit]
Description=Sync Google Drive every 15 minutes

[Timer]
OnBootSec=10m
OnUnitActiveSec=15m
Persistent=true

[Install]
WantedBy=timers.target
